# contrib/gh-dlr.timer
[Unit]
Description=Run gh-dlr weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
